# packaging/scripts/postinstall.sh
systemctl daemon-reload
systemctl enable prometheus-backfill-aws-athena.service
systemctl restart prometheus-backfill-aws-athena.service
