# packaging/scripts/preremove.sh
systemctl stop prometheus-backfill-aws-athena.service
systemctl disable prometheus-backfill-aws-athena.service
